# curator4-Mo.conf
[program:curator4-Mo]
command=python3 /home/curator4/curator4-Mo/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator4/curator4-Mo/test.err.log
stdout_logfile=/home/curator4/curator4-Mo/test.out.log